# os/petalinux_screen/project-spec/meta-user/recipes-apps/screen-boot-hook/files/screen-boot-hook.service
[Unit]
Description=Screen Boot Hook

[Service]
Type=simple
ExecStart=/usr/bin/screen-boot-hook.sh
Restart=no

[Install]
WantedBy=multi-user.target
